# ios/Library/Apple/System/Library/LaunchDaemons/com.apple.dt.RemotePairingDataVaultHelper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnablePressuredExit</key>
	<true/>
	<key>EnableTransactions</key>
	<true/>
	<key>Label</key>
	<string>com.apple.dt.RemotePairingDataVaultHelper</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.dt.RemotePairingDataVaultHelper</key>
		<true/>
	</dict>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>Program</key>
	<string>/Library/Apple/System/Library/PrivateFrameworks/RemotePairing.framework/Versions/A/Resources/bin/RemotePairingDataVaultHelper</string>
	<key>UserName</key>
	<string>_usbmuxd</string>
</dict>
</plist>
